# XG6846_4.12ITT01.83_20161115_consumer/userspace/private/apps/hotplug/hotplug.d/process/00wan
#!/bin/sh

[ "$NAME" = "wan" ] || exit
is_ipget="/var/run/ip.get"

if [ "$IPRENEWED" = "1" ]; then
	echo $WANPROTO > $is_ipget
fi

