# XG6846_4.12ITT01.83_20161115_consumer/userspace/private/apps/hotplug/hotplug.d/process/20httpd
#!/bin/sh

[ "$NAME" = "httpd" ] || exit

[ ! -s /var/authip ] || rm -f /var/authip
