# XG6846_4.12ITT01.83_20161115_consumer/userspace/private/apps/hotplug/hotplug.d/process/30trclient_acs
#!/bin/sh

[ "$NAME" = "trclient_acs" ] || exit
is_connected="/var/run/trclient_acs.connected"

if [ "$CONNECTED" = "1" ]; then
	echo > $is_connected
else
	rm -f $is_connected > /dev/null 2>&1
fi
