# XG6846_4.12ITT01.83_20161115_consumer/userspace/private/apps/hotplug/hotplug.d/net/00interfaces
#!/bin/sh

case "$ACTION" in
	add)
		;;
	*)
		;;
esac
